# systemd/path.path
[Unit]
Description=Process Incoming Files

[Path]
PathExistsGlob=/home/jfkelly/incoming/*



