# systemd/timer.timer
[Unit]
Description = Try Timer

[Timer]

OnBootSec=5
OnUnitActiveSec=5
AccuracySec=1

