# debugging/python/crack.service
[Unit]
Description=Password Cracking Service


[Service]
ExecStart=/root/presentations/debugging/python/crack.py

